# debian/sophomorix-base.cron.daily
#!/bin/sh
/usr/sbin/sophomorix-check --cron > /var/log/cron.daily.lastlog
