# debian/sophomorix-base.cron.weekly
#!/bin/sh

# log start
logger $0 started

# todo
echo 'ANALYZE' | psql -U ldap ldap && logger 'ANALYZE in postgres run succesfully' || logger 'ANALYZE in postgres run with ERROR'

# log end
logger $0 ended

