# debian/sophomorix-base.cron.monthly
#!/bin/sh

# log start
logger $0 started

# todo

# log end
logger $0 ended

